# scripts/crontab
SHELL=/bin/bash
0 3 * * * /usr/local/bin/python /app/install_smarttube.py >> /proc/1/fd/1 2>&1

